# deployment/kkl-backend.service
[Unit]
Description=KKL Forest Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kkl-forest/app_backend
Environment="PATH=/root/kkl-forest/app_backend/.venv/bin"
ExecStartPre=/bin/sh -c 'fuser -k 8000/tcp 2>/dev/null; sleep 1; true'
ExecStart=/root/kkl-forest/app_backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
